# resources/data/vpn-watchdog.service
[Unit]
Description=VPN Route Watchdog
After=network-online.target nss-lookup.target time-sync.target kodi.service
Wants=network-online.target nss-lookup.target time-sync.target
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/bin/python3 /storage/.kodi/addons/service.wireguard.manager/resources/lib/service.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
